# systemd/system/websudo-rootd.service
[Unit]
Description=websudo root executor
After=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/websudo/websudo.env
ExecStart=/usr/local/bin/websudo-rootd
Restart=on-failure

[Install]
WantedBy=multi-user.target
