# systemd/system/websudo-rootd@.service
[Unit]
Description=websudo root executor for uid %i
After=network.target

[Service]
Type=simple
Environment=WEBSUDO_ROOT_SOCKET_PATH=/run/user/%i/websudo-rootd.sock
Environment=WEBSUDO_ROOT_ALLOWED_UID=%i
EnvironmentFile=-/etc/websudo/websudo.env
ExecStart=/usr/local/bin/websudo-rootd
Restart=on-failure

[Install]
WantedBy=multi-user.target
